# openclaw-gateway.conf
[program:openclaw-gateway]
command=/usr/local/bin/start.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/openclaw-gateway-%(process_num)s.log
stderr_logfile=/var/log/openclaw-gateway-%(process_num)s.err
priority=999